# ibkr_bridge/systemd/ibkr-bridge.service
[Unit]
Description=IBKR Bridge API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ibkr
Group=ibkr
WorkingDirectory=/opt/stock-scanner-cloud
EnvironmentFile=/etc/ibkr-bridge.env
EnvironmentFile=-/etc/ibkr-gateway.env
ExecStart=/opt/stock-scanner-cloud/.venv/bin/python -m ibkr_bridge.app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
